# files/etc/supervisor/conf.d/run_tests.conf
[program:runtests]
command=/usr/local/bin/test_runner
autostart=true
autorestart=false
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
